# docker/supervisor.d/tmux.ini
[program:tmux]
command=/site/cli/docker tmux-start
environment =
    DB_HOST=%(ENV_DB_HOST)s,
    DB_PASSWORD=%(ENV_DB_PASSWORD)s,
    DB_PORT=%(ENV_DB_PORT)s,
    DB_USERNAME=%(ENV_DB_USERNAME)s,
    IS_DOCKER=%(ENV_IS_DOCKER)s,
    MYSQL_ROOT_PASSWORD=%(ENV_MYSQL_ROOT_PASSWORD)s,
    TZ=%(ENV_TZ)s,
    COLUMNS="207",
    LINES=46
user=nntmux
directory=/site
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_events_enabled=true
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_events_enabled=true
stderr_logfile_maxbytes=0
